# openrc/my-zram-config.in
#!/sbin/openrc-run

_command="my-zram-config"


depend() {
    after localmount
    before logger
}

start() {
    if ! "$_command" status; then
        "$_command" start
    fi
}

stop() {
    "$_command" stop
}
