# openrc/my-zram-config.in
#!/sbin/openrc-run

_command="my-zram-config"


depend() {
    after localmount
    before logger
}

start() {
    "$_command" status

    if [ "$?" != "0" ]; then
        "$_command" start
    fi
}

stop() {
    "$_command" stop
}
